# systemd/kasane-codex-watch.service
[Unit]
Description=kasane session watcher
After=default.target

[Service]
Type=simple
WorkingDirectory=%h/dev/kasane
ExecStart=%h/dev/kasane/.venv/bin/python -m kasane.main watch-all --codex-dir %h/.codex/sessions --opencode-db %h/.local/share/opencode/opencode.db --interval 30 --codex-limit 100 --opencode-limit 100
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
